# post_ip_to_slack

#! /bin/sh
### BEGIN INIT INFO
# Description:       This file should be used to post your ip address to slack on startup
#                    placed in /etc/init.d.
### END INIT INFO
# Author: Colin Kennedy

case "$1" in
  start)
        echo "Running" > /tmp/post_to_slack_log.txt
        ifconfig -a >> /tmp/post_to_slack_log.txt
        python /home/friday/python-scripts/post_ip_to_slack.py >> /tmp/post_to_slack_log.txt
        ifconfig -a >> /tmp/post_to_slack_log.txt
  ;;
  esac

  :
